# deploy/corretorcenter.service.example
[Unit]
Description=CorretorCenter
After=network.target postgresql.service

[Service]
Type=simple
WorkingDirectory=__WORKDIR__
EnvironmentFile=__WORKDIR__/.env
ExecStart=__NODE_BIN__ __WORKDIR__/app/server.js
Restart=always
RestartSec=5
User=__RUN_USER__
Group=__RUN_GROUP__

[Install]
WantedBy=multi-user.target
